# containerd/catalyst-redis.service
[Unit]
Description=Catalyst Redis Cache (containerd)
After=network-online.target containerd.service
Requires=containerd.service
Wants=network-online.target

[Service]
Type=forking
Restart=always
RestartSec=5s
TimeoutStartSec=120
TimeoutStopSec=30

# Environment variables
Environment="CONTAINER_NAME=catalyst-redis"
Environment="REDIS_IMAGE=docker.io/library/redis:7-alpine"
Environment="REDIS_PORT=6379"
Environment="DATA_DIR=/var/lib/catalyst/redis-data"
Environment="CTR_NAMESPACE=catalyst"

# Pre-start: ensure directory exists
ExecStartPre=-/usr/bin/mkdir -p /var/lib/catalyst/redis-data

# Main start command
ExecStart=/root/catalyst3/containerd/redis-wrapper-start.sh

# Stop command
ExecStop=/usr/bin/ctr -n catalyst task kill catalyst-redis SIGTERM
ExecStop=/bin/sleep 3
ExecStop=/usr/bin/ctr -n catalyst task rm -f catalyst-redis || true

# Post-stop
ExecStopPost=/bin/rm -f /tmp/catalyst-redis.pid

# Resource limits
LimitNOFILE=1048576
LimitNPROC=1048576

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
