# containerd/catalyst-seaweedfs.service
[Unit]
Description=Catalyst SeaweedFS Distributed Storage (containerd)
After=network-online.target containerd.service
Requires=containerd.service
Wants=network-online.target

[Service]
Type=forking
Restart=always
RestartSec=10s
TimeoutStartSec=180
TimeoutStopSec=30

# Environment variables
Environment="CONTAINER_NAME=catalyst-seaweedfs"
Environment="SEAWED_IMAGE=docker.io/chrislusf/seaweedfs:latest"
Environment="S3_PORT=8333"
Environment="MASTER_PORT=9333"
Environment="VOLUME_PORT=8080"
Environment="FILER_PORT=8888"
Environment="DATA_DIR=/var/lib/catalyst/seaweedfs"
Environment="CTR_NAMESPACE=catalyst"

# Pre-start: ensure directory exists
ExecStartPre=-/usr/bin/mkdir -p /var/lib/catalyst/seaweedfs

# Main start command
ExecStart=/root/catalyst3/containerd/seaweedfs-wrapper-start.sh

# Stop command
ExecStop=/usr/bin/ctr -n catalyst task kill catalyst-seaweedfs SIGTERM
ExecStop=/bin/sleep 5
ExecStop=/usr/bin/ctr -n catalyst task rm -f catalyst-seaweedfs || true

# Post-stop
ExecStopPost=/bin/rm -f /tmp/catalyst-seaweedfs.pid

# Resource limits
LimitNOFILE=1048576
LimitNPROC=1048576

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
