# containerd/catalyst-postgres.service
[Unit]
Description=Catalyst PostgreSQL Database (containerd)
After=network-online.target containerd.service
Requires=containerd.service
Wants=network-online.target

[Service]
Type=forking
Restart=always
RestartSec=5s
TimeoutStartSec=300
TimeoutStopSec=30

# Environment variables
Environment="CONTAINER_NAME=catalyst-postgres"
Environment="POSTGRES_IMAGE=docker.io/library/postgres:16-alpine"
Environment="POSTGRES_PORT=5432"
Environment="POSTGRES_USER=catalyst"
Environment="POSTGRES_PASSWORD=catalyst_dev_password"
Environment="POSTGRES_DB=catalyst_db"
Environment="DATA_DIR=/var/lib/catalyst/postgres-data"
Environment="CTR_NAMESPACE=catalyst-panel"

# Pre-start: ensure directory exists and has correct permissions
ExecStartPre=-/usr/bin/mkdir -p /var/lib/catalyst/postgres-data
ExecStartPre=/usr/bin/chown -R 999:999 /var/lib/catalyst/postgres-data

# Main start command - use the launcher script
ExecStart=/root/catalyst3/containerd/postgres-wrapper-start.sh

# Stop command - gracefully stop the container
ExecStop=/usr/bin/ctr -n catalyst task kill catalyst-postgres SIGTERM
ExecStop=/bin/sleep 5
ExecStop=/usr/bin/ctr -n catalyst task rm -f catalyst-postgres || true

# Post-stop - cleanup if needed
ExecStopPost=/bin/rm -f /tmp/catalyst-postgres.pid

# Resource limits
LimitNOFILE=1048576
LimitNPROC=1048576

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
