# containerd/catalyst-minio.service
[Unit]
Description=Catalyst MinIO Object Storage (containerd)
After=network-online.target containerd.service
Requires=containerd.service
Wants=network-online.target

[Service]
Type=forking
Restart=always
RestartSec=5s
TimeoutStartSec=180
TimeoutStopSec=30

# Environment variables
Environment="CONTAINER_NAME=catalyst-minio"
Environment="MINIO_IMAGE=docker.io/minio/minio:latest"
Environment="MINIO_PORT=9000"
Environment="CONSOLE_PORT=9001"
Environment="MINIO_ROOT_USER=minioadmin"
Environment="MINIO_ROOT_PASSWORD=minioadmin"
Environment="DATA_DIR=/var/lib/catalyst/minio"
Environment="CTR_NAMESPACE=catalyst"

# Pre-start: ensure directory exists
ExecStartPre=-/usr/bin/mkdir -p /var/lib/catalyst/minio

# Main start command
ExecStart=/root/catalyst3/containerd/minio-wrapper-start.sh

# Stop command
ExecStop=/usr/bin/ctr -n catalyst task kill catalyst-minio SIGTERM
ExecStop=/bin/sleep 5
ExecStop=/usr/bin/ctr -n catalyst task rm -f catalyst-minio || true

# Post-stop
ExecStopPost=/bin/rm -f /tmp/catalyst-minio.pid

# Resource limits
LimitNOFILE=1048576
LimitNPROC=1048576

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
